# config/s6/log/run
#!/bin/sh
# 这个脚本将实时显示 cf2dns.log 的内容
exec tail -f /tmp/cf2dns.log
