# config/s6/set-permissions/run
#!/bin/bash

# 启动赋予权限的脚本
exec /etc/s6/services/set-permissions/set_permissions.sh
